# launchd/com.shipscout.weeklyrefresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shipscout.weeklyrefresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-i</string>
    <string>/usr/local/bin/node</string>
    <string>/Users/macbook/Desktop/shipscout/scripts/weeklyRefresh.js</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/macbook/Desktop/shipscout</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/macbook/Desktop/shipscout/logs/weekly_refresh.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/macbook/Desktop/shipscout/logs/weekly_refresh.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
